# com.lex3.daily-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lex3.daily-news</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/xyp/.openclaw/workspace/scripts/daily-news.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/xyp/.openclaw/workspace/logs/daily-news.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xyp/.openclaw/workspace/logs/daily-news-error.log</string>
</dict>
</plist>